# config/systemd/storm-globus-gridftp.service
[Unit]
Description=StoRM Globus GridFTP service

[Service]
EnvironmentFile=-/etc/sysconfig/storm-globus-gridftp
User=root
Type=forking
ExecStart=/usr/sbin/globus-gridftp-server -S
KillMode=process

[Install]
WantedBy=multi-user.target
